# alfred.cron
# Cron

*/15 * * * * youruser cd /path/to/alfred/; ./env/bin/python ./build.py > ./feed.json
