# alfred.cron
# Cron

*/15 * * * * aaronb cd /home/aaronb/code/www/alfred/; ./env/bin/python ./build.py > ./feed.json
